# cpu_monitor.service
[Unit]
Description=CPU Monitor Daemon
After=network.target

[Service]
Type=simple
ExecStart=/home/nacho/bin/cpu_monitor.sh  # Replace with the actual path to your script
Restart=always

[Install]
WantedBy=multi-user.target

